# cron.d/syco-logviewer
*/5 * * * * root    /var/www/html/syco-signer/bin/mysql_clean_exclude.py
